# sensord/scripts/upstart.conf
description "sensord daemon"
author "Udi <udi@trugman.com>"

start on (started networking)
stop on runlevel [!023456]

respawn
kill signal SIGTERM

script
        exec cat
end script

